# meta-bah/recipes-application/bah/files/bah.service
[Unit]
Description=BAH

[Service]
ExecStart=/usr/bin/bah
Type=simple

[Install]
WantedBy=multi-user.target
